# configurations/ansible/roles/ngrok/templates/supervisor_ngrok.conf.j2
[program:ngrok]
command = /home/ubuntu/ngrok http -subdomain=crowdflare -config=/home/ubuntu/.ngrok2/ngrok.yml 8000
stdout_logfile=/var/log/ngrok_out.log
stderr_logfile=/var/log/ngrok_err.log
